# root/etc/services.d/rpcproxy/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc node /app/rpcproxy.js
